# infra/systemd/mysql.service.d/limits.conf.example
# MySQL のファイルディスクリプタ上限を引き上げる例
# 本番反映時は limits.conf にリネームして deploy してください
[Service]
LimitNOFILE=65535
